# ops/systemd/tttc-web-prod.service
[Unit]
Description=TTTC Web production (Vite preview)
After=network.target tttc-api-prod.service

[Service]
Type=simple
User=arved
Group=arved
WorkingDirectory=/srv/tttc/prod/web
ExecStart=/usr/bin/npm run preview:prod
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
# Include local node_modules/.bin so `vite` resolves (npm run should add this; systemd PATH alone is not enough if deps missing).
Environment=PATH=/srv/tttc/prod/web/node_modules/.bin:/usr/bin:/usr/local/bin

[Install]
WantedBy=multi-user.target
